# jenkins_install/files/jenkins.ini
[program:jenkins]
command = /data/apps/opt/jenkins/jenkins.sh
autostart = true
autorestart = true
priority = 30
startsecs = 5
startretries = 3
user = jenkins
redirect_stderr = false
stdout_logfile_backups= 5
stdout_logfile_maxbytes=200MB
stdout_events_enabled = false
stdout_capture_maxbytes=1MB
stderr_logfile_maxbytes=200MB
stderr_logfile_backups= 5
stderr_capture_maxbytes=1MB
stderr_events_enabled = false
stdout_logfile=/data/apps/log/jenkins/jenkins_std.log
stderr_logfile=/data/apps/log/jenkins/jenkins_err.log
