# ansible_runner_service_install/templates/ansible-runner-service.ini
[program:{{runner_service_app_name}}]
command = {{runner_service_install_dir}}/{{runner_service_app_name}}.sh
autostart = true
autorestart = true
priority = 30
startsecs = 5
startretries = 3
user = root
redirect_stderr = false
stdout_logfile_backups= 5
stdout_logfile_maxbytes=200MB
stdout_events_enabled = false
stdout_capture_maxbytes=1MB
#environment=ZOOCFGDIR="/etc/zookeeper"
stderr_logfile_maxbytes=200MB
stderr_logfile_backups= 5
stderr_capture_maxbytes=1MB
stderr_events_enabled = false
stdout_logfile={{runner_service_log_dir}}/{{runner_service_app_name}}_std.log
stderr_logfile={{runner_service_log_dir}}/{{runner_service_app_name}}_err.log
