# deploy_project/templates/project_common.ini
[program:{{project_prog_name}}]
command={%if project_prog_run_args is string%}{{project_prog_run_args}}{%else%}{{ project_prog_run_args | join(' ')}}{%endif%}

directory={{project_install_dir}}
#process_name=%(program_name)s_%(process_num)s
stdout_logfile={{project_log_dir}}/%(program_name)s_std.log
stdout_logfile_backups= 5
stdout_events_enabled = false
stdout_capture_maxbytes=1MB
stdout_logfile_maxbytes=200MB
stderr_logfile={{project_log_dir}}/%(program_name)s_err.log
stderr_logfile_maxbytes=200MB
stderr_logfile_backups= 5
stderr_capture_maxbytes=1MB
stderr_events_enabled = false
{%if project_prog_envs |length >0-%}
environment=
    {%-for env in project_prog_envs-%}
{{env}}={{project_prog_envs[env]}}
        {%-if env != project_prog_envs[-1]-%}
,
        {%-endif-%}
    {%-endfor-%}
{%-endif%}
numprocs=1
user = {{project_run_user}}
stopsignal=TERM
autostart=true
autorestart=true
