# logstash_install/templates/logstash.ini
[program:{{logstash_app_name}}]
#command = {{logstash_base_dir}}/bin/logstash -f {{logstash_conf_dir}}/%(program_name)s.conf -w 5 --path.data /data/apps/data/logstash/%(program_name)s --path.logs /data/apps/log/logstash/%(program_name)s
command = {{logstash_base_dir}}/bin/logstash --path.data {{logstash_data_dir}}/%(program_name)s --path.logs {{logstash_log_dir}}/%(program_name)s
directory = {{logstash_data_dir}}
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user = {{logstash_run_user}}
redirect_stderr = true
stdout_logfile={{logstash_log_dir}}/%(program_name)s_std.log
stderr_logfile={{logstash_log_dir}}/%(program_name)s_err.log
stdout_logfile_backups= 5
stdout_events_enabled = 1
stderr_logfile_backups= 5
stderr_events_enabled = 1
