# mysql_mha_install/templates/mha_boot.ini
[program:{{mha_instance_name}}]
command=/usr/bin/masterha_manager --conf={{ mha_manager_conf_dir }}/{{mha_instance_name}}.conf
directory={{mha_manager_log_dir}}/{{mha_instance_name}}
#process_name=%(program_name)s_%(process_num)s
numprocs=1
user = root
stopsignal=TERM
autostart=true
autorestart=true
stdout_logfile={{mha_manager_log_dir}}/%(program_name)s_std.log
stdout_logfile_backups= 5
stdout_events_enabled = 1
stderr_logfile={{mha_manager_log_dir}}/%(program_name)s_err.log
stderr_logfile_backups= 5
stderr_events_enabled = 1
