# filebeat_install/templates/filebeat.ini
[program:{{filebeat_app_name}}]
command = {{filebeat_base_dir}}/filebeat -c {{filebeat_conf_dir}}/filebeat.yml -e
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user = root
redirect_stderr = true
stdout_logfile={{filebeat_log_dir}}/filebeat-std.log
stderr_logfile={{filebeat_log_dir}}/filebeat-err.log
