# common_install/templates/elasticsearch_exporter.ini
[program:{{exporter_base_name}}]
command = {{exporter_base_dir}}/{{exporter_base_name}} --es.uri=http://{{ansible_default_ipv4.address}}:9200 --es.node="_local"  --log.level="info" --log.format="json" --log.output="stdout" --web.listen-address=":{{exporter_port}}"
autostart = true
autorestart = true
startsecs = 3
startretries = 3
user = {{exporter_run_user}}
redirect_stderr = false
stdout_logfile={{exporter_log_dir}}/%(program_name)s-std.log
stderr_logfile={{exporter_log_dir}}/%(program_name)s-err.log
