# logstash_conf/files/logstash_postloan.ini
[program:logstash_postloan3-sys]
command = /data/apps/opt/logstash/bin/logstash  -f /data/apps/config/logstash/logstash_postloan3-sys.conf -w 5 --path.data /data/apps/data/logstash/postloan3-sys --path.logs /data/apps/log/logstash/%(program_name)s --pipeline.id %(program_name)s
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user = es
redirect_stderr = true
stdout_logfile=/data/apps/log/logstash/%(program_name)s_std.log
stderr_logfile=/data/apps/log/logstash/%(program_name)s_err.log
