# logstash_conf/files/logstash_rcs2.ini
[program:logstash_rcs2-server-info]
command = /data/apps/opt/logstash/bin/logstash  -f /data/apps/config/logstash/logstash_rcs2-server-info.conf -w 5 --path.data /data/apps/data/logstash/rcs2-server-info --path.logs /data/apps/log/logstash/%(program_name)s --pipeline.id %(program_name)s
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user = es
redirect_stderr = true
stdout_logfile=/data/apps/log/logstash/%(program_name)s_std.log
stderr_logfile=/data/apps/log/logstash/%(program_name)s_err.log
[program:logstash_rcs2-scheduler-info]
command = /data/apps/opt/logstash/bin/logstash  -f /data/apps/config/logstash/logstash_rcs2-scheduler-info.conf -w 5 --path.data /data/apps/data/logstash/rcs2-scheduler-info --path.logs /data/apps/log/logstash/%(program_name)s --pipeline.id %(program_name)s
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user = es
redirect_stderr = true
stdout_logfile=/data/apps/log/logstash/%(program_name)s_std.log
stderr_logfile=/data/apps/log/logstash/%(program_name)s_err.log
[program:logstash_rcs2-jobcli-info]
command = /data/apps/opt/logstash/bin/logstash  -f /data/apps/config/logstash/logstash_rcs2-jobcli-info.conf -w 5 --path.data /data/apps/data/logstash/rcs2-jobcli-info --path.logs /data/apps/log/logstash/%(program_name)s --pipeline.id %(program_name)s
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user = es
redirect_stderr = true
stdout_logfile=/data/apps/log/logstash/%(program_name)s_std.log
stderr_logfile=/data/apps/log/logstash/%(program_name)s_err.log
[program:logstash_rcs2-api-worker-info]
command = /data/apps/opt/logstash/bin/logstash  -f /data/apps/config/logstash/logstash_rcs2-api-worker-info.conf -w 5 --path.data /data/apps/data/logstash/rcs2-api-worker-info --path.logs /data/apps/log/logstash/%(program_name)s --pipeline.id %(program_name)s
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user = es
redirect_stderr = true
stdout_logfile=/data/apps/log/logstash/%(program_name)s_std.log
stderr_logfile=/data/apps/log/logstash/%(program_name)s_err.log
[program:logstash_rcs2-abevent-worker-info]
command = /data/apps/opt/logstash/bin/logstash  -f /data/apps/config/logstash/logstash_rcs2-abevent-worker-info.conf -w 5 --path.data /data/apps/data/logstash/rcs2-abevent-worker-info --path.logs /data/apps/log/logstash/%(program_name)s --pipeline.id %(program_name)s
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user = es
redirect_stderr = true
stdout_logfile=/data/apps/log/logstash/%(program_name)s_std.log
stderr_logfile=/data/apps/log/logstash/%(program_name)s_err.log
