# prometheus_install/templates/prometheus.ini
[program:{{prom_base_name}}]
command = {{prom_base_dir}}/prometheus --config.file={{prom_conf_dir}}/prometheus.yml --web.listen-address="0.0.0.0:{{prom_port}}" --storage.tsdb.path={{prom_data_dir}} --web.enable-lifecycle
#directory= /data/apps/data/prometheus
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user = {{prom_run_user}}
redirect_stderr = false
stdout_logfile={{prom_log_dir}}/%(program_name)s_std.log
stdout_logfile_backups= 5
stdout_events_enabled = false
stdout_capture_maxbytes=1MB
stdout_logfile_maxbytes=200MB
stderr_logfile={{prom_log_dir}}/%(program_name)s_err.log
stderr_logfile_maxbytes=200MB
stderr_logfile_backups= 5
stderr_capture_maxbytes=1MB
stderr_events_enabled = false
