# redis_install/templates/redis.ini
[program:{{redis_service_name}}]
command={{redis_base_dir}}/bin/redis-server {{redis_conf_dir}}/{{redis_service_name}}.conf --daemonize no
directory={{redis_base_dir}}
user={{redis_run_user}}
priority=70
startsecs=10
autostart=true
autorestart=true
startretries=0
stopsignal=TERM
stopwaitsecs=5
redirect_stderr=false
stdout_logfile={{redis_log_dir}}/{{redis_service_name}}_std.log 
stdout_logfile_maxbytes=10MB  
stdout_logfile_backups=10    
stdout_capture_maxbytes=1MB  
stdout_events_enabled=false  
stderr_logfile={{redis_log_dir}}/{{redis_service_name}}_err.log
stderr_logfile_maxbytes=10MB   
stderr_logfile_backups=10     
stderr_capture_maxbytes=1MB  
stderr_events_enabled=false 
