# elasticsearch_install/templates/elasticsearch.ini
[program:elasticsearch]
command = {{es_base_dir}}/bin/elasticsearch
directory = {{es_base_dir}}
autostart = true
autorestart = true
startsecs = 10
startretries = 3
stopsignal=TERM
killasgroup = true
user = {{es_run_user}}
environment=JAVA_HOME="{{ansible_env.JAVA_HOME}}",PATH="{{ansible_env['PATH']}}"
redirect_stderr = false
stdout_logfile={{es_log_dir}}/%(program_name)s_std.log
stdout_logfile_maxbytkafka = 300MB
stdout_logfile_backups= 20
stdout_events_enabled = false
stderr_logfile={{es_log_dir}}/%(program_name)s_err.log
stderr_logfile_maxbytkafka=300MB
stderr_logfile_backups=10
stderr_capture_maxbytkafka=3MB
stderr_events_enabled= false
