# elasticsearch_install/templates/elasticsearch.ini
[program:elasticsearch]
command = {{es_base_dir}}/bin/elasticsearch
directory = {{es_base_dir}}
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user = {{es_run_user}}
environment=JAVA_HOME="{{ansible_env.JAVA_HOME}}"
redirect_stderr = false
stdout_logfile={{es_log_dir}}/elasticsearch-std.log
stderr_logfile={{es_log_dir}}/elasticsearch-err.log
