# flink_install/templates/flink_HistoryServer.ini
[program:flink_HistoryServer]
command={{flink_base_dir}}/bin/historyserver.sh start-foreground
directory={{flink_base_dir}}
stdout_logfile={{flink_log_dir}}/%(program_name)s-info.log
stdout_logfile_backups= 1
stdout_events_enabled = 1
stderr_logfile={{flink_log_dir}}/%(program_name)s-err.log
stderr_logfile_backups= 1
stderr_events_enabled = 1
environment=HDP_VERSION="$(hdp-select versions)",HADOOP_CONF_DIR="/usr/hdp/${HDP_VERSION}/hadoop/conf",HADOOP_CLASSPATH="$(hadoop classpath)"
numprocs=1
user = flink
stopsignal=TERM
autostart=true
autorestart=true

