# filebeat_install/files/filebeat.ini
[program:filebeat]
command = /data/apps/opt/filebeat/filebeat -c /data/apps/config/filebeat/filebeat.yml -e
autostart = true
autorestart = true
startsecs = 5
startretries = 3
user = root
redirect_stderr = true
stdout_logfile=/data/apps/log/filebeat/filebeat-std.log
stderr_logfile=/data/apps/log/filebeat/filebeat-err.log
